# autorestic.service
[Unit]
Description=Runs autorestic cron
Wants=autorestic.timer

[Service]
ExecStart=autorestic --ci cron -c /etc/autorestic/config.yml
